# config/supervisor/flask.conf
[program:flask]
command=python3 app.py
directory=/var/www/api
autostart=true
autorestart=true
startretries=3